# user/environment.target
[Unit]
Description=Set session environment variables
Requires=dbus.socket
IgnoreOnIsolate=true
