# user/environment.target
[Unit]
Description=Set session environment variables
IgnoreOnIsolate=true
