# user/wm.target
[Unit]
Description=Window manager
Wants=xorg.target
Wants=xinitrc.target
Requires=dbus.socket
AllowIsolate=true

[Install]
WantedBy=default.target
