# user/xdg-vars.service
[Unit]
Description=XDG variable setting

[Install]
WantedBy=environment.target
