# user/xinitrc.target
[Unit]
Description=Xinitrc Stuff
After=environment.target
Wants=multiplexer.target
#Wants=xfsettingsd.service

[Install]
Alias=default.target
