# user/xrdb.service
[Unit]
Description=X resource setting
Documentation=man:xrdb(1)

[Service]
Type=oneshot
ExecStart=/usr/bin/xrdb -merge %h/.Xresources

[Install]
WantedBy=xinitrc.target
