# user/xrdb.service
[Unit]
Description=X resource setting

[Service]
Type=oneshot
ExecStart=/usr/bin/xrdb -merge %h/.Xresources

[Install]
WantedBy=xinitrc.target
